# ixc_django_docker/etc/supervisord.nginx-proxy.tmpl.conf
[program:nginx]
autorestart = true
autostart = true
command = nginx.sh
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stopasgroup = true

[program:{{ default .Env.SUPERVISORD_PROGRAM "gunicorn" }}]
autorestart = true
autostart = true
command = {{ default .Env.SUPERVISORD_PROGRAM "gunicorn" }}.sh
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stopasgroup = true
